# task1_jboss/templates/wildfly.service.erb
[Unit]
Description=Wildfly
After=network.target

[Service]
Type=simple
User=<%= @user %>
Group=<%= @group %>
ExecStart=<%= @home %>/bin/standalone.sh 

[Install]
WantedBy=multi-user.target